# systemd/system/hermes-dashboard.service
[Unit]
Description=Hermes Agent Web Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/.hermes/hermes-agent
ExecStart=/home/ubuntu/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main dashboard --host 127.0.0.1 --port 9119 --no-open --insecure
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/.hermes/logs/dashboard.log
StandardError=append:/home/ubuntu/.hermes/logs/dashboard.log

[Install]
# Use `default.target` (not `multi-user.target`) so the unit can be
# enabled under either the system manager or the user manager.
# `multi-user.target` is system-only and breaks `enable` when the
# unit is symlinked into ~/.config/systemd/user/ by restore.sh
# (which is what happened on 2026-07-09 after PR #45 landed —
# restore.sh now symlinks every systemd/system/*.service into
# the user manager).
WantedBy=default.target
