# systemd/system/hermes-dashboard.service
[Unit]
Description=Hermes Agent Web Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/.hermes/hermes-agent
ExecStart=/home/ubuntu/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main dashboard --host 127.0.0.1 --port 9119 --no-open --insecure
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/.hermes/logs/dashboard.log
StandardError=append:/home/ubuntu/.hermes/logs/dashboard.log

[Install]
WantedBy=multi-user.target
