# systemd/system/mercury-state-backup.service
[Unit]
Description=Backup mercury host state (config + data) to /home/ubuntu/data/backups/
Documentation=file:///home/ubuntu/data/code/mercury-host/scripts/backup-mercury-state.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Refresh secrets dump before state backup so any state-side change that
# touched a .env lands in ~/.secrets/secrets.yaml within the same 24h
# window. --force overwrites in place (single-file, no retention).
# Idempotent: same script + same targets, daily cadence, no growth.
ExecStartPre=/home/ubuntu/data/code/mercury-host/scripts/backup-secrets.sh --force
ExecStart=/home/ubuntu/data/code/mercury-host/scripts/backup-mercury-state.sh
WorkingDirectory=/home/ubuntu
# systemd sets HOME=/ by default for system services; scripts under set -u
# explode on `${HOME:-$HOME/...}`. Pin HOME to the operator account so
# XDG_STATE_HOME defaults resolve correctly and tee-permission-sensitive
# paths land in a user-writable directory.
Environment="HOME=/home/ubuntu"
# Generous timeout — first backup with full state may take 1-3 min; integrity
# check adds another ~30s. 15min ceiling covers pathological cases.
TimeoutStartSec=900
# Surface exit code + last lines via journalctl
StandardOutput=journal
StandardError=journal
# Don't restart — backup is one-shot, retry would corrupt the manifest.
SuccessExitStatus=0