# systemd/user/discord-notify.service
[Unit]
Description=discord-notify — Discord merge-notification dispatcher. Receives HMAC-signed POSTs from webhook-server after successful PR merges, looks up the originating Discord thread, posts a formatted confirmation message.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/data/code/webhook-server
ExecStart=/home/ubuntu/.volta/tools/image/node/24.16.0/bin/node --experimental-strip-types /home/ubuntu/data/code/webhook-server/src/dispatch/server.ts
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Path: see webhook-server.service — must include user bin dirs for
# the same reasons (Discord bot token lookup, etc.).
Environment=PATH=/home/ubuntu/.volta/bin:/home/ubuntu/.local/bin:/home/ubuntu/.local/share/pnpm/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/ubuntu
Environment=USER=ubuntu

# Discord bot token. Loaded from ~/.hermes/.env at setup time. Required.
# Set this to the bot token for the guild/channel where notifications
# should land. Use `hermes gateway setup` to bootstrap, or paste directly.
Environment=DISCORD_BOT_TOKEN=

# Override port/bind via env if needed (defaults below).
Environment=DISCORD_NOTIFY_PORT=3940
Environment=DISCORD_NOTIFY_BIND=127.0.0.1

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/ubuntu/.config/discord-notify /home/ubuntu/.hermes

[Install]
WantedBy=default.target
