# systemd/user/webhook-server.service
[Unit]
Description=webhook-server — GitHub webhook receiver for OpenSpec arc auto-merge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/data/code/webhook-server
ExecStart=/home/ubuntu/.volta/tools/image/node/24.16.0/bin/node --experimental-strip-types /home/ubuntu/data/code/webhook-server/src/server.ts
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Path: extend the systemd default to include user bin dirs where pnpm,
# volta-shim, and project-local tools live. The default user-service
# PATH is just /usr/{local/,}bin which doesn't include ~/.local/bin or
# ~/.volta/bin. Without this, `exec('pnpm', ...)` fails with ENOENT.
#
# This mirrors what login shells get from ~/.profile:
#   export PATH="$VOLTA_HOME/bin:$PATH"
Environment=PATH=/home/ubuntu/.volta/bin:/home/ubuntu/.local/bin:/home/ubuntu/.local/share/pnpm/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/ubuntu
Environment=USER=ubuntu

# Hardening: don't run as root if possible. We're already 'ubuntu'.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/ubuntu/data/.webhook-server /home/ubuntu/.config/webhook-server /home/ubuntu/data/code /home/ubuntu/data/code/webhook-server /home/ubuntu/data/code/x-digest /home/ubuntu/data/code/better-bet

[Install]
WantedBy=default.target
